# deploy/digitalocean/liquidator-scheduler.service
[Unit]
Description=Liquidator Articles Scheduler (one slot)
After=network-online.target liquidator-bot.service
Wants=network-online.target

[Service]
Type=oneshot
User=appuser
Group=appuser
WorkingDirectory=/home/appuser/apps/liquidator
EnvironmentFile=/home/appuser/apps/liquidator/.env
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
# Гарантируем что критичные python-зависимости в venv до старта слота.
# Если хоть один импорт упал — автоматически переустанавливаем requirements.
# Защита от ситуации когда venv пересоздан без полной переустановки
# (10 мая 2026 fal-client пропал → статьи сутки выходили без обложек).
# Накладные ~200мс при нормальной ситуации; ~30-90с при self-heal.
ExecStartPre=/bin/bash -c '/home/appuser/apps/liquidator/.venv/bin/python -c "import fal_client, cloudinary, httpx, dotenv" || /home/appuser/apps/liquidator/.venv/bin/pip install -r /home/appuser/apps/liquidator/requirements.txt'
ExecStart=/home/appuser/apps/liquidator/.venv/bin/python -m articles_scheduler.runner

Nice=10
# TimeoutStartSec = SLOT_BUDGET_SEC (5100) + запас на git push/finalize/self-heal
# pip install. На VPS уже стояло 6000 — синхронизируем шаблон с реальностью.
TimeoutStartSec=6000
StandardOutput=journal
StandardError=journal
SyslogIdentifier=liquidator-scheduler

NoNewPrivileges=true
ProtectSystem=full
ProtectHome=read-only
# Claude CLI пишет session state в /home/appuser/.claude.json (файл в HOME).
# Без него /write-article падает с EROFS read-only file system.
ReadWritePaths=/home/appuser/apps/liquidator /home/appuser/.git-credentials /home/appuser/.claude /home/appuser/.claude.json
PrivateTmp=true

# Claude CLI + контекст агентов up to 1.5-1.8 ГБ. Запас.
MemoryMax=2000M
